# modules/provision_ronin/templates/puppet-apply.timer.epp
[Unit]
Description=Run puppet-apply every <%= @puppet_run_interval %>

[Timer]
OnActiveSec=1m
OnUnitActiveSec=<%= @puppet_run_interval %>

[Install]
WantedBy=timers.target
